# patches/etc/rc.local
/usr/local/bin/add_htpasswd_user.sh
/usr/local/bin/configure_ssl.sh

exit 0
